# src/system-configs/upstart-configs/units/webida-proxy.conf
#!upstart
description "Webida proxy"
author "S-Core"

start on startup
stop on shutdown

expect fork

env WEBIDA_HOME="/home/webida"
env SERVER_DIR="webida-server/src/server"
env LOG="/dev/null"
env NODE_ENV='production'

env SCRIPT="unit-manager.js svc=proxy0"

script
	exec forever --sourceDir "$WEBIDA_HOME/$SERVER_DIR" -a -l $LOG --minUptime 5000 --spinSleepTime 2000 start $SCRIPT
end script

pre-stop script
	exec forever stop $SCRIPT
end script

