# src/server/buildjm/dist/webida-build-jm.conf
#!upstart
description "Webida Build Job Manager"
author "S-Core"

start on startup
stop on shutdown

# Set the priv of scripts
# If this not not set, the Webida server will be run as root and it'll be a disaster.
setuid webida

# This let upstart track the child process of forever.
# Because forever spawns another monitor process and exits, this is required.
expect fork

env DIR="/home/webida/server-build-jm"
env LOG="/var/webida/log/server-build-jm.log"
env SCRIPT="server-build-jm.js"

# Set HOME because forever keeps its setting(pids, socks,..) in $HOME/.forever
# If this is not set, /root/.forever will be used and fail to write.
env HOME="/home/webida"

chdir /home/webida/server-build-jm
script
	export HOME="/home/webida"
	export JAVA_HOME="/usr/lib/jvm/java-7-oracle"
	export ANDROID_HOME="$HOME/adt/sdk/tools"
	export ANDROID_PLATFORM_TOOLS="$HOME/adt/sdk/platform-tools"
	export PATH="$ANDROID_HOME:$ANDROID_PLATFORM_TOOLS:$PATH"
	export ANT_HOME="/usr/share/ant"
	export PATH="$PATH:$ANT_HOME/bin"
	exec forever --sourceDir $DIR -a -l $LOG --minUptime 5000 --spinSleepTime 2000 start $SCRIPT
end script

pre-stop script
	exec forever stop $SCRIPT
end script

