# src/server/buildjm/dist/webida-build.conf
#!upstart
description "Webida Build Server"
author "S-Core"

start on startup
stop on shutdown

# Set the priv of scripts
# If this not not set, the Webida server will be run as root and it'll be a disaster.
setuid webida

# This let upstart track the child process of forever.
# Because forever spawns another monitor process and exits, this is required.
expect fork

env DIR="/home/webida/server-build"
env LOG="/var/webida/log/server-build.log"
env SCRIPT="server-build.js"

# Set HOME because forever keeps its setting(pids, socks,..) in $HOME/.forever
# If this is not set, /root/.forever will be used and fail to write.
env HOME="/home/webida"

chdir /home/webida/server-build
script
	exec forever --sourceDir $DIR -a -l $LOG --minUptime 5000 --spinSleepTime 2000 start $SCRIPT
end script

pre-stop script
	exec forever stop $SCRIPT
end script

