# dist/target/default/etc/init/webida-proxy.conf
#!upstart
description "Webida proxy"
author "S-Core"

start on startup
stop on shutdown

expect fork

env SERVER_ROOT="/home/webida/webida-server/default"
env SERVER_DIR="server"
env LOG="/var/webida/log/proxy.log"
env SCRIPT="unit-manager.js svc=proxy0"

script
	exec forever --sourceDir "$SERVER_ROOT/$SERVER_DIR" -a --minUptime 5000 --spinSleepTime 2000 start $SCRIPT
end script

pre-stop script
	exec forever stop $SCRIPT
end script

