# dist/target/default/etc/init/webida.conf
#!upstart
description "Webida Server"
author "S-Core"

start on startup
stop on shutdown

# Set the priv of scripts
# If this not not set, the Webida server will be run as root and it'll be a disaster.
setuid webida

# This let upstart track the child process of forever.
# Because forever spawns another monitor process and exits, this is required.
expect fork

env SERVER_ROOT="/home/webida/webida-server/default"
env SERVER_DIR="server"
env LOG="/var/webida/log/app.log"
env SCRIPT="unit-manager.js svc=app0"

# Set HOME because forever keeps its setting(pids, socks,..) in $HOME/.forever
# If this is not set, /root/.forever will be used and fail to write.
env HOME="/home/webida"

script
	exec forever --sourceDir "$SERVER_ROOT/$SERVER_DIR" -a --minUptime 5000 --spinSleepTime 2000 start $SCRIPT
end script

pre-stop script
	exec forever stop $SCRIPT
end script

